# devops/debian/prerm
#!/bin/sh
# prerm script for osconfig
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <prerm> `remove'
#        * <old-prerm> `upgrade' <new-version>
#        * <new-prerm> `failed-upgrade' <old-version>
#        * <conflictor's-prerm> `remove' `in-favour' <package> <new-version>
#        * <deconfigured's-prerm> `deconfigure' `in-favour'
#          <package-being-installed> <version> `removing'
#          <conflicting-package> <version>
# for details, see https://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    remove|upgrade|deconfigure)
    ;;

    failed-upgrade)
    ;;

    *)
        echo "prerm called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

# Automatically added by dh_installsystemd/12.10ubuntu1
if [ -d /run/systemd/system ] && [ "$1" = remove ]; then
    # MSFT change: dh_installsystemd uses deb-systemd-invoke which fails
    # when /usr/sbin/policy-rc.d fails, use systemctl directly
    systemctl stop 'osconfig.service' >/dev/null || true
    systemctl disable 'osconfig.service' >/dev/null || true
    systemctl stop 'osconfig-platform.service' >/dev/null || true
    systemctl disable 'osconfig-platform.service' >/dev/null || true
fi
# End automatically added section

# MSFT change: Remove module identity on package removal
if [ "$1" = remove ]; then
    rm /etc/aziot/identityd/config.d/osconfig.toml >/dev/null || true
fi

exit 0
